# temporary/testDaemon.service
[Unit]
Description=Test Daemon Service

[Service]
Type=forking
ExecStart=/usr/bin/daemonProc
 
[Install]
WantedBy=multi-user.target
